# src/bin/patchmanager-daemon/systemd/system-update-patchmanager-pre.service
[Unit]
Description=Move away parts of Patchmanager which might affect a system update
Before=system-update.target

[Service]
Type=oneshot
ExecStart=
ExecStart=-/usr/bin/cp /etc/ld.so.preload /etc/ld.so.preload_pm3.bak
ExecStart=s/usr/bin/ed -i '/libpreloadpatchmanager/ d' /etc/ld.so.preload
ExecStartPost=-/sbin/ldconfig

[Install]
Wants=system-update-pre.target
